# conf/mt/crontab
*/3 * * * * cd /var/www/cgi-bin/mt; ./tools/run-periodic-tasks
